# packaging/kpulse-daemon.user.service
[Unit]
Description=KPulse system heartbeat daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/kpulse-daemon
Restart=on-failure

[Install]
WantedBy=default.target
